# launchd/com.cass.semantic-reindex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cass.semantic-reindex</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>launchctl unload ~/Library/LaunchAgents/com.cass.index-watch.plist 2>/dev/null; sleep 2; /Users/__USERNAME__/.local/bin/cass index --semantic; launchctl load ~/Library/LaunchAgents/com.cass.index-watch.plist</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/__USERNAME__/Library/Logs/cass-semantic-reindex.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USERNAME__/Library/Logs/cass-semantic-reindex.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/__USERNAME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/__USERNAME__</string>
    </dict>
</dict>
</plist>
